# finnScan.conf
[program:finn-timer-app]
directory=/home/soumya/python/finnScan
command=python timer.py
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
